# deploy/systemd/buttonsbebe-whatsapp-connect.service.d/password.conf
[Service]
Environment=WA_PASSWORD=<REDACTED>
