# deploy/systemd/buttonsbebe-processor.service.d/whatsapp.conf
[Service]
Environment="WHATSAPP_SEND_URL=<REDACTED>"
