# deploy/systemd/buttonsbebe-webhook.service
[Unit]
Description=Buttons Bebe webhook receiver (Gorgias -> job queue)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/Buttonsbebe Agent/webhook
Environment="PYTHONPATH=/root/Buttonsbebe Agent/webhook/src"
ExecStart=/bin/bash -c 'exec "/root/Buttonsbebe Agent/webhook/.venv/bin/uvicorn" bb_webhook.app:app --host 127.0.0.1 --port 8000'
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
