# gorgias-webhook/gorgias-webhook.service
[Unit]
Description=Gorgias Webhook Receiver
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/gorgias-webhook
ExecStart=/usr/bin/python3 /root/gorgias-webhook/server.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target