# gorgias-webhook/infra/hindsight/hindsight.service
[Unit]
Description=Hindsight Agent Memory Service (Docker container)
After=docker.service
Requires=docker.service

[Service]
Type=simple
ExecStart=/usr/bin/docker start -a hindsight
ExecStop=/usr/bin/docker stop hindsight
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target